# monitoring/watch-farmer.timer
[Unit]
Description=Watch Farm Timer

[Timer]
OnActiveSec=30
OnUnitActiveSec=30

[Install]
WantedBy=timers.target