# harvest/start-harvester.service
[Unit]
After=network.service

[Service]
ExecStart=/root/chia-blockchain/venv/bin/chia start harvester

[Install]
WantedBy=default.target